# target/services/init-resize/init-resize.service
[Unit]
Description=Init Resize
 
[Service]
ExecStart=/usr/local/bin/init-resize.sh
 
[Install]
WantedBy=multi-user.target
